# brigsd.service
[Unit]
Description=Brigs: Backlight Daemon

[Service]
Type=forking
ExecStart=/vault/code/brigs/brigsd.sh

[Install]
WantedBy=multi-user.target
